# systemd/clusterchk.socket
[Unit]
Description=Clusterchk socket

[Socket]
ListenStream=10.10.24.164:9999
Accept=true

[Install]
WantedBy=sockets.target
