# systemd/clusterchk.socket
[Unit]
Description=Clusterchk socket

[Socket]
ListenStream=9999
Accept=true

[Install]
WantedBy=sockets.target
